# body/services/nox-body.service
[Unit]
Description=Nox Body Controller (PiDog)
After=network.target

[Service]
EnvironmentFile=/home/pidog/nox.env
Type=simple
User=pidog
WorkingDirectory=/home/pidog
ExecStart=/usr/bin/python3 -u /home/pidog/nox_daemon.py
Restart=on-failure
RestartSec=5
Environment=PATH=/home/pidog/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=SDL_AUDIODRIVER=alsa
# AUDIODEV intentionally not set here: the daemon auto-detects the speaker;
# pin it in nox.env instead if detection picks the wrong card.

[Install]
WantedBy=multi-user.target
